# root/etc/services.d/rclone-mount/run
#!/usr/bin/with-contenv sh

umask 022

if [ -z "$REMOTE_MOUNT" ]; then
  echo "Error: REMOTE_MOUNT environment variable was not passed to the container."
  exit 1
fi

if [ -z "$LOCAL_MOUNT" ]; then
  echo "Error: LOCAL_MOUNT environment variable was not passed to the container."
  exit 1
fi

if [ -z "$RCLONE_MOUNT_OPTIONS" ]; then
  echo "Error: RCLONE_MOUNT_OPTIONS environment variable was not passed to the container."
  exit 1
fi

mount_command="rclone mount $RCLONE_MOUNT_OPTIONS $REMOTE_MOUNT $LOCAL_MOUNT"

echo "Executing => $mount_command"
exec \
	s6-setuidgid abc $mount_command